# crontab
#run python script every minutes
* * * * * python /app/foxit.py >> run.log 2>&1 &